# service/openai_api_status_prober.service
[Unit]
Description=Run the node.js application for openai-api-status-prober
After=network.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/usr/local/bin/openai_status_prober_app  # Symlink to your app directory
ExecStart=node /usr/local/bin/openai_status_prober_app/server.js  # Path to the Node.js executable and your app's main file
Restart=on-failure
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=openai_status_prober

[Install]
WantedBy=multi-user.target